# home/.config/systemd/user/dump-packages.timer
[Install]
WantedBy=timers.target

[Timer]
OnCalendar=*:0/1

[Unit]
Description=Timer for dump-packages service
